# maintain/build/llama.cpp-service/llama.cpp.socket
[Unit]
Description=Llama.cpp Socket
After=network.target

[Socket]
ListenStream=/run/llama.cpp.sock

[Install]
WantedBy=sockets.target
